# scripts/kongtrade-wallet-report.service
[Unit]
Description=KongTradeBot Weekly Wallet Scout Report
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/KongTradeBot
ExecStart=/usr/bin/python3 /root/KongTradeBot/scripts/weekly_wallet_report.py
EnvironmentFile=/root/KongTradeBot/.env
StandardOutput=journal
StandardError=journal
