# scripts/kongtrade-skipped-report.timer
[Unit]
Description=KongTradeBot Weekly Skipped-Signal Report Timer

[Timer]
OnCalendar=Fri *-*-* 23:55:00 Europe/Berlin
Persistent=true

[Install]
WantedBy=timers.target
