# scripts/kongtrade-signal-evaluator.timer
[Unit]
Description=KongTradeBot Daily Skipped-Signal Evaluator Timer

[Timer]
OnCalendar=*-*-* 02:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
