# scripts/kongtrade-signal-evaluator.service
[Unit]
Description=KongTradeBot Daily Skipped-Signal Evaluator
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/claudeuser/KongTradeBot
ExecStart=/usr/bin/python3 /home/claudeuser/KongTradeBot/scripts/run_signal_evaluator.py
EnvironmentFile=/home/claudeuser/KongTradeBot/.env
StandardOutput=journal
StandardError=journal
