# scripts/kongtrade-wallet-report.service
[Unit]
Description=KongTradeBot Weekly Wallet Scout Report
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/claudeuser/KongTradeBot
ExecStart=/usr/bin/python3 /home/claudeuser/KongTradeBot/scripts/weekly_wallet_report.py
EnvironmentFile=/home/claudeuser/KongTradeBot/.env
StandardOutput=journal
StandardError=journal
