# scripts/kongtrade-tax-export.service
[Unit]
Description=KongTradeBot Weekly Tax Export
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/claudeuser/KongTradeBot
ExecStart=/usr/bin/python3 /home/claudeuser/KongTradeBot/scripts/weekly_tax_export.py
EnvironmentFile=/home/claudeuser/KongTradeBot/.env
StandardOutput=journal
StandardError=journal
